# ansible/templates/cv-builder.service.j2
[Unit]
Description=CV Builder API (Node.js)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory={{ cv_app_dir }}/server
EnvironmentFile={{ cv_app_dir }}/server/.env
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
